# pkg/dom0-ztools/rootfs/etc/profile.d/path.sh
# make PATH variable stable
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
export PATH
